# conf/beanstalkd.conf
[program:beanstalkd]
command = /usr/bin/beanstalkd -l 0.0.0.0 -p 11300 -b /var/lib/beanstalkd
user = beanstalkd
autostart=true
autorestart=true
